# timeshift-support/grub-btrfs-timeshift.path
[Unit]
Description=Monitors for new snapshots

[Path]
PathModified=/run/timeshift/backup/timeshift-btrfs/snapshots

[Install]
WantedBy=multi-user.target
